# dev-setup/backend/etc/systemd/system/pip_install_api.path
[Unit]
Description=Monitor /var/www/backend/api/pip.requirements for changes.
Wants=pip_install_api.service

[Path]
PathChanged=/var/www/backend/api/pip.requirements

[Install]
WantedBy=multi-user.target
